# backend/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=/bin/bash start.sh

[program:celery]
command=celery -A core worker --loglevel=info